# crontab
*/5 * * * * /var/www/html/cacti/scripts/cacti-vsphere.py vcsa