# crontab
*/5 * * * * /usr/bin/python /var/www/html/cacti/scripts/cacti-vsphere.py vcenterhostname